# repo/meta-apollo/recipes-core/systemd-apollo/files/data.mount
[Mount]
What=/dev/mmcblk0p3
Where=/data
Options=noatime
